# Network/Ansible/ansible_linux_update_final/service_config/run_local.service
[Unit]
Description=Run_Local.sh
After=network-online.target

[Service]
ExecStart=/home/aniket/run_local.sh

[Install]
WantedBy=multi-user.target
